# ansible/roles/web/templates/status_monitor.service.j2
[Unit]
Description=Cloud Pulse Web Status Monitor
After=network.target

[Service]
Type=simple
User=root
ExecStart=/usr/bin/python3 /opt/monitoring/status_monitor.py
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
